# ProyectoF/Aplicacion/supervisord.conf
[supervisord]
nodaemon=true

[program:mysql]
command=/etc/init.d/mysql start &
autostart=true
autorestart=true
user=root

[program:apache]
command=apache2ctl -DFOREGROUND